# .config/systemd/user/check-pw-expiration.service
[Unit]
Description=Check for and alert about approaching local password expiration

[Service]
Type=oneshot
ExecStart=%h/.df/bin/check-pw-expiration
